# Config/com.clashpow.engine.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.clashpow.engine</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Applications/ClashPow.app/Contents/MacOS/clashpow-engine</string>
    </array>

    <key>MachServices</key>
    <dict>
        <key>com.clashpow.engine</key>
        <true/>
    </dict>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
        <key>Crashed</key>
        <true/>
    </dict>

    <key>ProcessType</key>
    <string>Background</string>

    <key>LowPriorityIO</key>
    <true/>

    <!-- Crash recovery: restart after 3 seconds on failure -->
    <key>ThrottleInterval</key>
    <integer>3</integer>

    <key>EnvironmentVariables</key>
    <dict>
        <key>GOMAXPROCS</key>
        <string>4</string>
        <key>GOGC</key>
        <string>50</string>
    </dict>

    <key>StandardOutPath</key>
    <string>/usr/local/var/log/clashpow-engine.log</string>

    <key>StandardErrorPath</key>
    <string>/usr/local/var/log/clashpow-engine.log</string>
</dict>
</plist>
